# build/scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:emulator]
command=/usr/local/bin/run-emulator.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/emulator.stdout.log
stderr_logfile=/var/log/emulator.stderr.log
priority=5
stopsignal=TERM
stopwaitsecs=50
stopasgroup=true
killasgroup=true

[program:socat_console]
command=/usr/bin/socat TCP-LISTEN:5554,fork,reuseaddr TCP:127.0.0.1:5558
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/socat_console.stdout.log
stderr_logfile=/var/log/socat_console.stderr.log
priority=10
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:socat_adb]
command=/usr/bin/socat TCP-LISTEN:5555,fork,reuseaddr TCP:127.0.0.1:5559
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/socat_adb.stdout.log
stderr_logfile=/var/log/socat_adb.stderr.log
priority=10
stopsignal=TERM
stopasgroup=true
killasgroup=true


[program:adb_frida]
command=/usr/local/bin/adb-frida-forward.sh
autostart=true
autorestart=true
startretries=3
priority=12
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/adb_frida.stdout.log
stderr_logfile=/var/log/adb_frida.stderr.log

; optional env overrides
environment=SERIAL="emulator-5558",REMOTE_PORT="27042",LOCAL_PORT="37042"


[program:socat_frida]
command=/usr/bin/socat TCP-LISTEN:37043,fork,reuseaddr TCP:127.0.0.1:37042
autostart=true
autorestart=true
startretries=3
priority=12
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/socat_frida.stdout.log
stderr_logfile=/var/log/socat_frida.stderr.log